# gtw.service
[Unit]
Description=Gateway IoT
After=systemd-networkd.service

[Service]
ExecStart= python /root/gateway.py

[Install]
WantedBy=multi-user.target
